# Web/Web-summit-works/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
directory=/src
command=./server.sh

[program:client]
directory=/src
command=./client.sh
